# jetson/.bin/resources/supervisor/conf.d/indiweb.conf
[program:indi-web]
;directory=%(ENV_CONDA_PREFIX)s/bin
directory=/opt/miniforge/bin
command=indi-web
autorestart=true
stopasgroup=true
killasgroup=true
priority=200